# ansible/roles/linux-base/templates/etc/profile.d/prompt-color.sh.j2
#!/bin/bash
export PS1="{{ prompt[0] }}"


