# ansible/roles/provisioning/templates/etc/profile.d/rbenv.sh
export RBENV_ROOT="/usr/lib/rbenv"
export PATH="/usr/lib/rbenv/bin:$PATH"
eval "$(rbenv init -)"
